# condor_data_collectors/csjobs.service
[Unit]
Description=Queries the condor Schedd daemon for job info and sends it off to redis
After=network.target

[Service]
Type=simple
User=csjobs
ExecStart=/usr/bin/python /opt/cloudscheduler/condor_data_collectors/csjobs.py
Restart=on-abort


[Install]
WantedBy=multi-user.target
