# etc/cron.d/csv2-copy-certificates-to-users.afile
# run once a week
0 0 * * Mon root /usr/local/sbin/copy_letsencrypt_to_users.sh
