# etc/cron.d/ec2_retrieve_flavor_files
SHELL=/bin/bash
PATH=/bin:/usr/bin:/usr/local/bin:/sbin:/usr/sbin:/usr/local/sbin
MAILTO=
0 0 * * * cloudscheduler /opt/cloudscheduler/utilities/ec2_retrieve_flavor_files.py
