# etc/cron.d/generate-public-page.afile
# regenerate public static page every 10 minutes
*/10 * * * * python3 /opt/cloudscheduler/web_frontend/cloudscheduler/csv2/gen_public_page.py
