# etc/systemd/system/csv2-openstack.service
[Unit]
Description= Retrieves a list of openstack clouds from the configured database and queries them for metadata which is relayed back to the database
After=network.target

[Service]
Type=simple
User=cloudscheduler
ExecStart=/usr/bin/python3 /opt/cloudscheduler/data_collectors/openstack/openstackPoller.py
Restart=always

[Install]
WantedBy=multi-user.target
